# extras/supervisor/flask.conf
[program:flask]
user=ubuntu
command=/usr/local/bin/flask run
directory=/home/ubuntu/shared/microblog

