# xport/extras/supervisor/flask.conf
[program:xport]
user=ubuntu
command=/usr/local/bin/xport-run
