# xport/extras/supervisor/flask.conf
[program:xport]
command=/usr/local/bin/xport-run
